# meson/meson.service
[Unit]
Description="Meson systemd service"
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=meson_service
SyslogIdentifier=meson
SyslogFacility=local7
KillSignal=SIGHUP
EnvironmentFile=/etc/meson/meson-service.env
ExecStart=/home/b6111427/meson-linux-amd64/meson
WorkingDirectory=/home/b6111427/meson-linux-amd64
[Install]
WantedBy=multi-user.target
